# crontab
0 2 * * * /path/to/sync_leetcode_problems.sh

